# borg/backup.timer


[Unit]
Description=Run borg backup on /home/matteo folder weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target