# borg/backup.service
[Unit]
Description=create home backup service


[Service]
ExecStart=borg create /storage/borg::Home-$(date -I) /home/matteo
