# borg/backup.service
[Unit]
Description=create home backup service


[Service]
ExecStart=/home/matteo/scripts/borg/home_backup.sh 
