# timers/backup.timer


[Unit]
Description=Run borg backup on user home folder weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target