# timers/grive.service
[Unit]
Description=sync google drive with drive


[Service]
ExecStart=/home/matteo/scripts/grive/sync.sh
