# timers/weather.service
[Unit]
Description=run weather notifications 


[Service]
ExecStart=/home/matteo/scripts/weather/weather.sh 
