# systemd/user/timers/backup.service
[Unit]
Description=create home backup 


[Service]
ExecStart=/home/matteo/scripts/systemd/scripts/home_backup.sh 
