# systemd/user/timers/grive.service
[Unit]
Description=sync google drive with drive


[Service]
ExecStart=/home/matteo/scripts/systemd/scripts/sync.sh
