# systemd/user/timers/rclone-mount-all.timer


[Unit]
Description=Run rclone mount all script

[Timer]
OnStartupSec=3seconds
Persistent=true

[Install]
WantedBy=timers.target
