# systemd/timers/backup.service
[Unit]
Description=Backup your files with borg

[Service]
ExecStart=[[HOME]]/scripts/bin/backup.sh backup
