# systemd/timers/backup.service
[Unit]
Description=Run home_backup service


[Service]
ExecStart=/home/matteo/scripts/systemd/scripts/home_backup.sh 
