# systemd/timers/backup.timer


[Unit]
Description=Run backup service

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
