# systemd/timers/backup.service
[Unit]
Description=Run home_backup service

[Service]
ExecStart=[[HOME]]/scripts/systemd/bin/backup.sh backup
