# systemd/timers/git_sync.service
[Unit]
Description=Sync git repositories

[Service]
ExecStart=[[HOME]]/scripts/systemd/bin/git_sync.sh sync
