# systemd/timers/newsboat.service
[Unit]
Description=updates newsboat feed

[Service]
ExecStart=[[HOME]]/scripts/systemd/bin/newsboat.sh
