# systemd/timers/prune.service
[Unit]
Description=prune borg repo periodically

[Service]
ExecStart=[[HOME]]/scripts/systemd/bin/backup.sh prune
