# systemd/timers/rclone-mount-all.service
[Unit]
Description=mount rclone storage 


[Service]
ExecStart=/home/matteo/scripts/systemd/scripts/rclone_mount_all.sh 
