# services/backup.service
[Unit]
Description=create home backup 


[Service]
ExecStart=/home/matteo/scripts/borg/home_backup.sh 
