# systemd/vpn.timer
[Unit]
Description=Timer for the vpn service
Requires=vpn.service

[Timer]
OnBootSec=180sec

[Install]
WantedBy=timers.target
