# systemd/wallpaper_randomizer.service
[Unit]
Description=Change current wallpaper

[Service]
ExecStart=/bin/bash ${HOME}/.local/bin/wallpaper_randomizer.sh change_wall
