# systemd/.config/systemd/user/backup.service
[Unit]
Description=Backup your files with borg

[Service]
ExecStart=/bin/bash ${HOME}/.local/bin/backup.sh backup
