# systemd/.config/systemd/user/newsboat.service
[Unit]
Description=updates newsboat feed

[Service]
ExecStart=/bin/bash ${HOME}/.local/bin/newsboat.sh
